# .config/autostart/keyboard.desktop
[Desktop Entry]
Name=Keyboard
Exec=$KEYBOARD_CONFIG
Terminal=false
Type=Application
